# systemd/inventory-web.service
[Unit]
Description=Inventory Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=inventory
Group=inventory
WorkingDirectory=/home/tobias/furusetalle9/inventory
ExecStart=/usr/bin/inventory-system serve
Restart=on-failure
RestartSec=5s

# Security
NoNewPrivileges=true
PrivateTmp=true
ReadOnlyPaths=/home/tobias/furusetalle9/inventory

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inventory-web

[Install]
WantedBy=multi-user.target
