# systemd/inventory-chat.service
[Unit]
Description=Inventory Chatbot Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=inventory
Group=inventory
WorkingDirectory=/home/tobias/furusetalle9/inventory
Environment="ANTHROPIC_API_KEY="
ExecStart=/usr/bin/inventory-system chat
Restart=on-failure
RestartSec=5s

# Security
NoNewPrivileges=true
PrivateTmp=true
ReadOnlyPaths=/home/tobias/furusetalle9/inventory

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inventory-chat

[Install]
WantedBy=multi-user.target
