# systemd/inventory-chat@.service
[Unit]
Description=Inventory Chatbot Server (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=inventory-%i
Group=inventory-%i
EnvironmentFile=/etc/inventory-system/%i.conf
Environment="ANTHROPIC_API_KEY=${ANTHROPIC_API_KEY}"
ExecStart=/opt/inventory-system/venv/bin/inventory-system chat ${INVENTORY_PATH} --port ${CHAT_PORT}
Restart=on-failure
RestartSec=5s

# Security
NoNewPrivileges=true
PrivateTmp=true
ReadOnlyPaths=${INVENTORY_PATH}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inventory-chat-%i

[Install]
WantedBy=multi-user.target
